# nifstd/resources/ontree.socket
[Unit]
Description=ontree gunicorn socket

[Socket]
ListenStream=/run/ontree/socket

[Install]
WantedBy=sockets.target
